# gramps-mute.service
[Unit]
Description=Gramps Phone Mute Helper
After=sound.target

[Service]
Type=notify
Environment=PYTHONUNBUFFERED=1
ExecStart=%h/gramps-env/bin/python %h/gramps-transcriber/mute_helper.py
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
