# systemd/caption.service
[Unit]
Description=Gramps Caption Service (BULLETPROOF v2.0)
After=sound.target graphical-session.target
StartLimitIntervalSec=300
StartLimitBurst=10

[Service]
Type=simple
Environment=DISPLAY=:0
Environment=WAYLAND_DISPLAY=wayland-0
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=PYTHONUNBUFFERED=1
ExecStartPre=/bin/sleep 3
ExecStart=%h/gramps-env/bin/python %h/gramps-transcriber/caption_app.py
StandardOutput=append:%h/caption.log
StandardError=append:%h/caption.log
Restart=always
RestartSec=5
WatchdogSec=120
TimeoutStartSec=300
TimeoutStopSec=10

# Resource limits
MemoryMax=1G
CPUQuota=90%

[Install]
WantedBy=default.target
